# velbusd/rootfs/etc/s6-overlay/s6-rc.d/velbusd/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Runs velbusd
# ==============================================================================
bashio::log.info "Starting VelbusD..."
exec python3 -m velbustcp --settings /settings.json
