# velbusd/rootfs/etc/s6-overlay/s6-rc.d/velbusd/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Runs velbusd
# ==============================================================================
declare configuration_file

configuration_file="/settings.json"

bashio::log.info "Starting VelbusD..."
exec python3 -m velbustcp --settings $(configuration_file)
