# velbusd/rootfs/etc/s6-overlay/s6-rc.d/velbusd/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Runs velbusd
# ==============================================================================

bashio::log.info "Building configfile..."
jq -n '{
  "ntp": {
    "enabled": $ntp,
    "synctime": "03:00"
  },
  "connections": [
    {
      "host": "0.0.0.0",
      "port": 27015,
      "relay": true,
      "ssl": false,
      "cert": "",
      "pk": "",
      "auth": false,
      "auth_key": ""
    }
  ],
  "serial": {
    "autodiscover": false,
    "port": ""
  },
  "logging": {
    "type": $logging,
    "output": "stream"
  }
}' \
  --arg ntp $(bashio::config 'ntp') \
  --arg logging $(bashio::config 'logging') \
> /settings.json


# starting
bashio::log.info "Starting VelbusD..."
exec python3 -m velbus --settings /settings.json
