# velbusd/rootfs/etc/s6-overlay/s6-rc.d/velbusd/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Runs velbusd
# ==============================================================================

bashio::log.info "Building configfile..."
bashio::log.info "ntp $(bashio::config 'ntp')"
bashio::log.info "logging $(bashio::config 'logging')"
bashio::log.info "velbus $(bashio::config 'velbus')"
bashio::log.info "port $(bashio::addon.port 27015)"
bashio::log.info "auth $(bashio::config 'auth')"
bashio::log.info "authkey $(bashio::config 'authkey')"
jq -n '{
  "ntp": {
    "enabled": $ntp | test("true"),
    "synctime": "03:00"
  },
  "connections": [
    {
      "host": "0.0.0.0",
      "port": $port | tonumber,
      "relay": true,
      "ssl": false,
      "cert": "",
      "pk": "",
      "auth": $auth | test("true"),
      "auth_key": $authkey
    }
  ],
  "serial": {
    "autodiscover": false,
    "port": $velbus
  },
  "logging": {
    "type": $logging,
    "output": "stream"
  }
}' \
  --arg ntp $(bashio::config 'ntp') \
  --arg logging $(bashio::config 'logging') \
  --arg velbus $(bashio::config 'velbus') \
  --arg port $(bashio::addon.port 27015) \
  --arg auth $(bashio::config 'auth') \
  --arg authkey "$(bashio::config 'authkey')" \
> /settings.json

cat /settings.json



# starting
bashio::log.info "Starting VelbusD..."
exec python3 -m velbustcp --settings /settings.json
