# velbusd/rootfs/etc/s6-overlay/s6-rc.d/velbusd/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Runs velbusd
# ==============================================================================

# set velbus port
bashio::config.require 'velbus'
jq ".serial.autodiscover=false" /settings.json > /tmp/set
cp /tmp/set /settings.json
jq ".serial.port=$(bashio::config 'velbus')" /settings.json > /tmp/set
cp /tmp/set /settings.json

# set log level
bashio::config.require 'logging'
jq ".logging.type=$(bashio::config 'logging')" /settings.json > /tmp/set
cp /tmp/set /settings.json

# enable/disable ntp
bashio::config.require 'ntp'
if bashio::config.true 'ntp'; then
     jq ".ntp.enabled=true" /settings.json > /tmp/set
else
     jq ".ntp.enabled=false" /settings.json > /tmp/set
fi
cp /tmp/set /settings.json

bashio::log.info "Starting VelbusD..."
exec python3 -m velbustcp --settings /settings.json
